# rc.local
#!/bin/sh -e
# run script below at boot:
/root/safe.sh
exit 0


